# scripts/docker-tcp.socket
[Unit]
Description=Docker Socket for the API

[Socket]
ListenStream=2385
BindIPv6Only=both
Service=docker.service

[Install]
WantedBy=sockets.target
